# neo4j_simple/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid
logfile=/var/log/supervisor/supervisord.log

[program:neo4j]
command="/srv/neo4j/bin/neo4j start"
